# root/etc/services.d/htpcmanager/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid mab python /app/htpcmanager/Htpc.py --datadir /config
